# goon.ini
[program:gostar]
directory= /data/svr/gostar
command = /data/svr/gostar/gostar -path /data/config/gostar.toml
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true